# programs/sophiatxd/package/debian/postinst
#!/bin/sh
# postinst script for sophiatx
#
# see: dh_installdeb(1)

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
# <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
# <failed-install-package> <version> `removing'
# <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
        DATA_DIR=/var/lib/sophiatx
        LOG_DIR=/var/log/sophiatx
        LOG_FILE=$LOG_DIR/sophiatx.log
        ENV_DEFAULT_FILE=/etc/default/sophiatx

        CONFIG_DIR=/etc/sophiatx
        SYSLOG_CONFIG_FILE=/etc/rsyslog.d/22-sophiatx.conf
        LOGROTATE_CONFIG_FILE=/etc/logrotate.d/sophiatx

        # Check if systemctl is present
        command -v systemctl >/dev/null 2>&1
        if [ "$?" -ne 0 ]; then
            echo "System must support systemctl to run sophiatx." >&2
            exit 1
        fi

        # Create user sophiatx if does not exist
        id sophiatx >/dev/null 2>&1
        if [ "$?" -ne 0 ]; then
            useradd --system -U -M sophiatx -s /bin/false -d $DATA_DIR
        fi

        test -d $DATA_DIR || mkdir -p $DATA_DIR
        test -d $LOG_DIR || mkdir -p $LOG_DIR
        touch $LOG_FILE

        # Adjust files permissions and owners
        test -d $CONFIG_DIR && chmod 644 $CONFIG_DIR/*
        test -f $SYSLOG_CONFIG_FILE && chmod 644 $SYSLOG_CONFIG_FILE
        test -f $LOGROTATE_CONFIG_FILE && chmod 644 $LOGROTATE_CONFIG_FILE

        chown -R -L sophiatx:sophiatx $DATA_DIR
        chown -R -L syslog:adm $LOG_DIR
        chmod 755 $DATA_DIR
        chmod 755 $LOG_DIR
        chmod 644 $LOG_FILE

        # Add defaults file, if it doesn't exist
        test -f $ENV_DEFAULT_FILE || touch $ENV_DEFAULT_FILE

        systemctl enable sophiatx
     ;;

     abort-upgrade|abort-remove|abort-deconfigure)
     ;;

     *)
         echo "postinst called with unknown argument \`$1'" >&2
         exit 1
     ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0